# .config/systemd/etc/systemd/user/pcmanfm.service
[Unit]
Description=Automatically mounts removable media
Documentation=https://wiki.archlinux.org/title/PCManFM#Daemon_mode
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
ExecStart=/bin/sh -c 'pcmanfm-qt -d || pcmanfm -d'
NonBlocking=true
ExecReload=kill -SIGUSR2 $MAINPID
Restart=on-failure

[Install]
WantedBy=graphical-session.target
